# .config/systemd/user/duplicitybak.timer
[Unit]
Description=Run duplicity backup every 33 days

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target